# .config/systemd/user/aur-check-updates.service
[Unit]
Description=Check AUR updates via auracle

[Service]
Type=simple
ExecStart=/usr/bin/bash /home/vk/apps/aur-check-updates.bash
