# .config/systemd/user/aur-check-updates.timer
[Unit]
Description=Checks AUR for updates every 6 hours

[Timer]
OnBootSec=1min
OnUnitActiveSec=6h

[Install]
WantedBy=graphical-session.target
