# .config/systemd/user/battery_low.timer
[Unit]
Description=Check battery life every so often

[Timer]
OnUnitActiveSec=10min

[Install]
WantedBy=default.target
