# .config/systemd/user/battery_low.timer
[Unit]
Description=Check battery life every so often

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=graphical-session.target
