# .config/systemd/user/battery_low.service
[Unit]
Description=Check battery life

[Service]
Type=simple
ExecStart=/home/vk/apps/battery.bash
