# .config/systemd/user/bt-stop-music.service
[Unit]
Description=Stop music when headphones disconnect

[Service]
Type=simple
ExecStart=/home/vk/apps/bt-stop-music.bash

[Install]
WantedBy=graphical-session.target
