# .config/systemd/user/nm-applet.service
[Unit]
Description=Network Manager applet

[Service]
Type=simple
ExecStart=/usr/bin/nm-applet

[Install]
WantedBy=xorg.target
