# .config/systemd/user/numlockx.service
[Unit]
Description=Enable numlock

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/numlockx on

[Install]
WantedBy=xorg.target
