# .config/systemd/user/thinkpad-keys.service
[Unit]
Description=Thinkpad key translator

[Service]
Type=simple
ExecStart=/home/vk/apps/thinkpad-keys.bash

[Install]
WantedBy=xorg.target
