# .config/systemd/user/thinkpad_keys.service
[Unit]
Description=Thinkpad key translator

[Service]
Type=simple
ExecStart=/home/vk/apps/thinkpad_keys.bash

[Install]
WantedBy=default.target
