# .config/systemd/user/xbindkeys.service
[Unit]
Description=xbindkeys service

[Service]
Type=simple
ExecStart=/usr/bin/xbindkeys -n

[Install]
WantedBy=default.target
