# .config/systemd/user/xorg.target
[Unit]
Description=Xorg running
BindsTo=graphical-session.target
